# mattermost-newsfeeds.service
[Unit]
Description=mattermost-newsfeeds Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/mattermost-newsfeeds/.venv/bin/python /opt/mattermost-newsfeeds/src/main.py
WorkingDirectory=/opt/mattermost-newsfeeds
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
